# mpv_mods/Etc_Systemd_user/vita.service
[Service]
ExecStart=/usr/bin/vlc v4l2:///dev/video0 :v4l2-standard= :live-caching=0 :v4l2-fps=60.00