# mpv_mods/Etc_Systemd_user/vita.service
[Service]
ExecStart=/usr/bin/mpv /dev/video0 --fullscreen --profile=low-latency --fps=60  --framedrop=no --speed=1.21 --really-quiet 